# images/proton-bridge/services/metrics/run
#!/bin/sh
set -eu

state_dir="/tmp/proton-bridge-s6"
mkdir -p "${state_dir}"
date +%s > "${state_dir}/metrics.start"

web_root="/tmp/proton-bridge-metrics-www"
mkdir -p "${web_root}/cgi-bin"
ln -sf /app/metrics/metrics.cgi "${web_root}/cgi-bin/metrics"

exec busybox httpd -f -p "${CONTAINER_METRICS_PORT}" -h "${web_root}"
