# images/proton-bridge/services/exporter/finish
#!/bin/sh
set -eu

exec /app/services/.s6-finish-policy.sh "exporter" "false" "${1:-0}" "${2:-0}"
